# nox-mem/nox-mem-watcher.service
[Unit]
Description=NOX-Supermem File Watcher — monitora memory/ e dispara ingest automático
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash __NOX_MEM_PATH__/nox-mem-watch.sh
Restart=always
RestartSec=5
StandardOutput=append:/tmp/nox-mem-watcher.log
StandardError=append:/tmp/nox-mem-watcher.log

[Install]
WantedBy=multi-user.target
